# bmac3-collector.service
[Unit]
Description=BMAC3 SWIM Flight Data Collector
After=network.target

[Service]
Type=simple
User=bmacdonald3
WorkingDirectory=/home/bmacdonald3
ExecStart=/usr/bin/python3 /home/bmacdonald3/collector.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
